# nkn/roles/nkn-systemd/templates/nkn.service
[Unit]
Description=Nkn service

[Service]
Restart=always
ExecStart=/bin/bash /nkn/entrypoint.sh

[Install]
WantedBy=multi-user.target
