# nkn/roles/nkn-azure/files/nkn.service
[Unit]
Description=Nkn service

[Service]
ExecStart=/bin/bash /nkn/entrypoint.sh

[Install]
WantedBy=multi-user.target
